# proxy/squid-supervisor.conf
[program:squid]
command=/usr/sbin/squid -N -d 1
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/squid.log
user=root
priority=1

[program:blacklist_watchdog]
command=/usr/bin/python3 /usr/local/bin/blacklist_watchdog.py
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/squid/watchdog.log
stdout_logfile_maxbytes=1MB
user=root
priority=2