# proxy/squid-supervisor.conf
[program:squid]
command=/usr/sbin/squid -N -d 1
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/squid.log
user=root
priority=1